# ansible/templates/prefect-worker-cpu.service.j2
[Unit]
Description=Prefect CPU worker
After=network.target
Requires=network.target

[Service]
Type=simple
User=prefect
Group=prefect
WorkingDirectory=/opt/prefect
Environment=PREFECT_HOME=/opt/prefect/.prefect
Environment=PREFECT_API_URL=http://{{ prefect_server_ip }}:4201/api
Environment=PYTHONPATH=/opt/prefect/flows

# identifier used for namespacing gpu ids
Environment=WORKER_ID={{ inventory_hostname }}

ExecStart=/opt/prefect/pixi-env/env/bin/prefect worker start --pool cpu-pool --limit {{ cpu_worker_limit | default(8) }}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
