# ansible/templates/prefect-worker.service.j2
[Unit]
Description=Prefect worker
After=network.target
Requires=network.target

[Service]
Type=exec
User=prefect
Group=prefect
WorkingDirectory=/opt/prefect
Environment=PREFECT_HOME=/opt/prefect/.prefect
Environment=PREFECT_API_URL=http://{{ prefect_server_ip }}:4201/api
Environment=PYTHONPATH=/opt/prefect/flows

{% if worker_gpu_ids is defined and worker_gpu_ids | length > 0 %}
# comma-separated list of GPU IDs, e.g. "0,1,2"
Environment=WORKER_GPU_IDS={{ worker_gpu_ids }}
{% endif %}

# identifier used for namespacing gpu ids
Environment=WORKER_ID={{ inventory_hostname }}

ExecStart=/opt/prefect/pixi-env/env/bin/prefect worker start --pool {{ work_pool_name }} --limit {{ worker_limit | default(8) }}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
