# ansible/templates/prefect-worker-gpu.service.j2
[Unit]
Description=Prefect GPU worker
After=network.target
Requires=network.target

[Service]
Type=simple
User={{ prefect_user }}
Group={{ prefect_user }}
WorkingDirectory={{ prefect_home }}
Environment=PREFECT_HOME={{ prefect_home }}/.prefect
Environment=PREFECT_API_URL=http://{{ prefect_server_ip }}:4201/api
Environment=PYTHONPATH={{ prefect_home }}/flows

# comma-separated list of GPU IDs, e.g. "0,1,2"
Environment=WORKER_GPU_IDS={{ worker_gpu_ids }}

# identifier used for namespacing gpu ids
Environment=WORKER_ID={{ inventory_hostname }}

ExecStart={{ prefect_home }}/pixi-env/env/bin/prefect worker start --pool gpu-pool --limit {{ gpu_worker_limit | default(8) }}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
